# bin/smilecdr.service
# ----------------------------------------------------------------------------
# Smile CDR  //  https://smilecdr.com
#
# Copyright 2019 Smile CDR, Inc.
# ----------------------------------------------------------------------------
#
# Keep smiling, this is the systemd init script for Linux hosts.
#
# For instructions on installing/using this script, please see
# https://smilecdr.com/docs/installation/unix_service_installation.html
#
# ----------------------------------------------------------------------------

[Unit]
Description=Smile CDR

[Service]
ExecStart=/opt/smilecdr/bin/smilecdr run
User=smile

Type=simple
LimitNOFILE=500000
LimitNPROC=500000
KillSignal=SIGINT
TimeoutStopSec=20
Restart=always
RestartSec=5
StartLimitInterval=0
StartLimitBurst=10
StandardOutput=null
StandardError=syslog

[Install]
WantedBy=multi-user.target
